# start_swift.service
[Unit]
Description=Mount Swift Drives by Zulfikar

[Service]
ExecStart=/opt/swift/bin/mount_devices.sh

[Install]
WantedBy=multi-user.target
